# postgresql/root/etc/services.d/adempiere/run
#!/usr/bin/with-contenv sh

cd /data/app/Adempiere/utils
./RUN_Server2.sh
